# deploy/sentinel-manager.service
[Unit]
Description=Valkey Sentinel Manager
Documentation=https://github.com/chals-go/valkey-sentinel-manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sentinel-manager
Group=sentinel-manager
ExecStart=/usr/local/bin/sentinel-manager --config /etc/sentinel-manager/config.yaml
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/sentinel-manager /etc/sentinel-manager
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-manager

[Install]
WantedBy=multi-user.target
